# scripts/pilot_guardian.service.example
[Unit]
Description=AgentLens Pilot Guardian
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/agentlens
Environment=AGENTLENS_BASE_URL=https://www.agentlens.one
Environment=PILOT_ALERT_WEBHOOK_URL=
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/agentlens/scripts/pilot_guardian.py --api-key <pilot_key> --customer "Bibin" --interval-sec 300
Restart=always
RestartSec=15
User=agentlens
Group=agentlens

[Install]
WantedBy=multi-user.target

